# home/profile.d/gtags.sh
# gtags (GNU Global) configuration
# Loaded by ~/.profile via profile.d mechanism

# Use Python pygments for better language support
export GTAGSLABEL='native-pygments'
export GTAGSCONF="$HOME/.global/gtags.conf"
